# monroe-experiments/DEBIAN/postinst
#!/bin/sh

systemctl enable monroe-namespace || true
systemctl restart monroe-namespace || true
