# monroe-experiments/DEBIAN/prerm
#!/bin/sh

systemctl stop monroe-namespace
systemctl disable monroe-namespace
